# misc/csengo-log.service
[Service]
Type=simple
ExecStart=qterminal -e /home/studio/log.sh

[Install]
WantedBy=default.target
