# misc/csengo.service
[Service]
Type=simple
ExecStart=/home/studio/csengo/csengo

[Install]
WantedBy=default.target
