# misc/csengo.service
[Service]
Type=simple
ExecStart=/home/studio/csengo/csengo
EnvironmentFile=/home/studio/csengo/env

[Install]
WantedBy=default.target
